# dracut/99setup-root/initrd-setup-root.service
[Unit]
Description=Root filesystem setup
DefaultDependencies=no
RequiresMountsFor=/sysroot/usr/
RequiresMountsFor=/sysroot/usr/share/oem
After=sysroot-usr.mount
After=initrd-root-fs.target
Before=initrd-parse-etc.service

[Service]
Type=oneshot
ExecStart=/sbin/initrd-setup-root
